# devbox-disk/etc/update-motd.d/10-devbox-login
#!/bin/bash
 echo;
    echo You are now logged in to;
    figlet $(hostname);
    cat  <<EOF

Mapped docker ports: $(cat /etc/docker.ports)

Host c-drive is mounted to /mnt/c/
EOF