# devbox-disk/etc/profile.d/.bash_aliases
alias idea=/usr/share/intellij/bin/idea.sh
alias ideaq='idea > /dev/null 2>&1 &'
